# scripts/debian/nucleoid.service
[Unit]
Description=Nucleoid Runtime

[Service]
ExecStart=/usr/bin/npm start --prefix /usr/lib/nucleoid/

[Install]
WantedBy=multi-user.target
