# scripts/debian/postinst
#!/bin/sh
systemctl start nucleoid
systemctl enable nucleoid
