# scripts/debian/prerm
#!/bin/sh
systemctl stop nucleoid
systemctl disable nucleoid
